# backend/fastmcp-seo-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.robofy.fastmcp-seo-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mathewpm/dev/BMAD-Method/robofy/backend/manage-fastmcp-server.sh</string>
        <string>start-foreground</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/fastmcp-seo-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/fastmcp-seo-server.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SERPAPI_API_KEY</key>
        <string>3d3aa8879cfc13c2bd19d8ded8b5e8afed9922b26d7b6473d370bf9cbaa194dd</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mathewpm/dev/BMAD-Method/robofy/backend</string>
</dict>
</plist>